# netplex/files/netplex.init
#!/bin/sh /etc/rc.common
START=99
STOP=10
USE_PROCD=1
. /usr/share/pcdn/netplex_update
start_service(){
    mk_pear_dir
    procd_open_instance netplex
    procd_set_param command /usr/bin/netplex
    procd_close_instance
}

stop_service(){
    echo stop
}
