# luci-app-fogvdn/root/etc/uci-defaults/luci-fogvdn
#!/bin/sh
rm -f /tmp/luci-indexcache
#uci track
cfg_id=$(uci add ucitrack fogvdn)
uci set ucitrack.$cfg_id.init="fogvdn" 2>/dev/null
uci commit ucitrack 2>/dev/null
exit 0
